# services/rustify-embed/rustify-embed.service
[Unit]
Description=rustify-embed (MERT audio embedding service for rustify-player)
# systemd --user cannot reference system units (docker.service is system-level).
# User is responsible for ensuring Docker is running before starting this unit.

[Service]
Type=simple
Restart=on-failure
RestartSec=5

# Fully-qualified image name locally. Build with:
#   cd /home/opc/rustify-player/services/rustify-embed
#   docker build -t rustify-embed:latest .
ExecStartPre=-/usr/bin/docker stop rustify-embed
ExecStartPre=-/usr/bin/docker rm rustify-embed
ExecStart=/usr/bin/docker run \
    --name rustify-embed \
    --rm \
    -p 127.0.0.1:8448:8448 \
    --cpus=14 \
    --memory=8g \
    -v rustify-embed-models:/models \
    rustify-embed:latest

ExecStop=/usr/bin/docker stop rustify-embed

[Install]
WantedBy=default.target
